# panzaiko.service
[Unit]
Description=Panzaiko Inventory Management System
After=network.target

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/var/www/panzaiko
Environment="PATH=/var/www/panzaiko/venv/bin"
Environment="FLASK_ENV=production"
ExecStart=/var/www/panzaiko/venv/bin/gunicorn --config gunicorn.conf.py app.backend.app:app
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
